# systemd/user/chromium.slice
[Unit]
Description=Chromium slice

[Slice]
MemoryAccounting=true
MemoryHigh=6G
