# systemd/user/chromium.slice
[Unit]
Description=Chromium slice

[Slice]
MemoryAccounting=true
MemoryLimit=6G
CPUAccounting=true
CPUQuota=50%
