# systemd/user/dlnar@.service
[Unit]
Description=dlnar
After=network.target

[Service]
ExecStart=/home/joram/bin/dlnar %i

[Install]
WantedBy=default.target
