# systemd/user/dlnar@.service
[Unit]
Description=dlnar
After=network.target

[Service]
ExecStart=%h/bin/dlnar %i

[Install]
WantedBy=default.target
