# systemd/user/firefox.slice
[Unit]
Description=Firefox slice

[Slice]
MemoryAccounting=true
MemoryHigh=8G
