# systemd/user/firefox.slice
[Unit]
Description=Firefox slice

[Slice]
MemoryAccounting=true
MemoryLimit=4G
CPUAccounting=true
CPUQuota=50%
