# systemd/user/httpr@.service
[Unit]
Description=mountr
After=network.target

[Service]
ExecStart=/home/joram/bin/httpr %i

[Install]
WantedBy=default.target
