# systemd/user/mountr@.service
[Unit]
Description=mountr
After=network.target

[Service]
ExecStart=/home/joram/bin/mountr %i

[Install]
WantedBy=default.target
