# systemd/user/x11vnc.service
[Unit]
Description=VNC server
After=network.target display-manager.service

[Service]
ExecStart=/usr/bin/x11vnc -forever -display WAIT:0.0
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
